# dsi_env/resize.service
[Unit]
Description=Service to resize SD card

[Service]
Type=oneshot
ExecStart=/bin/bash /root/resize

[Install]
WantedBy=multi-user.target
